# deploy/packages/emqx.service
[Unit]
Description=emqx daemon
After=network.target

[Service]
User=emqx
Group=emqx
Type=forking
Environment=HOME=/var/lib/emqx

# Must use a 'bash' wrap for some OS
# errno=13 'Permission denied'
# Cannot create FIFO ... for writing
ExecStart=/bin/bash /usr/bin/emqx start

LimitNOFILE=1048576
ExecStop=/bin/bash /usr/bin/emqx stop
Restart=on-failure

# When clustered, give the peers enough time to get this node's 'DOWN' event
RestartSec=60s

[Install]
WantedBy=multi-user.target
